# jhub/jupyterhub.service
[Unit]
Description=Jupyter Hub Server

[Service]
#Type=forking
PIDFile=/var/run/jupyterhub.pid

# Jupyter Notebook: change PATHs as needed for your system
ExecStart=/usr/local/bin/jupyterhub -f /etc/jupyterhub/jupyterhub_config.py
# ExecStop=/bin/kill -s TERM $MAINPID
StandardOutput=file:/etc/jupyterhub/main.log
StandardError=file:/etc/jupyterhub/error.log


User=root
# Group=root
WorkingDirectory=/etc/jupyterhub
# KillMode=none
# Restart=on-failure

[Install]
WantedBy=multi-user.target
